# crontab
0 * * * * /host/job/bin/run-monitor.sh -e
